# meta-amd/meta-amd-bsp/recipes-core/udev/udev-extraconf/60-drm-hotplug-mode.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", ENV{DISPLAY}=":0", RUN+="/etc/udev/scripts/hotplug-display-@MULTIDISPLAYMODE.sh"
